# crontab
# START CRON JOB
0 0 * * * cd ../app && python3 scraper.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB
